# containers/rootfs/etc/s6-overlay/s6-rc.d/lskypro/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck -d
/usr/local/sbin/php /var/www/lsky-pro/artisan queue:work --queue=emails,images,thumbnails