# containers/rootfs/etc/cron.d/lskypro
* * * * * root /usr/local/sbin/php /var/www/lsky-pro/artisan schedule:run